# core/federation-dashboard.service
[Unit]
Description=S30-PRO Federation Dashboard (External Link)
After=network.target mnt-raid1.mount

[Service]
Type=simple
User=ed
WorkingDirectory=/home/ed/seevar/core/dashboard
ExecStart=/home/ed/.pyenv/versions/ssc-3.13.5/bin/python3 dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
