# systemd/seeing-scraper.service
[Unit]
Description=Meteoblue Seeing Scraper
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
ExecStart=%h/seevar/.venv/bin/python3 core/preflight/seeing_scraper.py

TimeoutStartSec=30s
TimeoutStopSec=10s
