# systemd/seevar-orchestrator.service
[Unit]
Description=SeeVar Science Orchestrator
After=network-online.target seevar-weather.service seevar-gps.service
Wants=network-online.target seevar-weather.service seevar-gps.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/seevar/.venv/bin/python3 core/flight/orchestrator.py
Restart=on-failure
RestartSec=15
StandardOutput=append:%h/seevar/logs/orchestrator.log
StandardError=append:%h/seevar/logs/orchestrator.log

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ReadWritePaths=%h/seevar/data %h/seevar/logs
TimeoutStartSec=30s
TimeoutStopSec=10s

[Install]
WantedBy=default.target
