# systemd/seevar-planner.service
[Unit]
Description=SeeVar Nightly Planner
After=network-online.target seevar-weather.service seevar-gps.service
Wants=network-online.target seevar-weather.service seevar-gps.service

[Service]
Type=oneshot
WorkingDirectory=%h/seevar
Environment=PYTHONPATH=%h/seevar
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/seevar/.venv/bin/python3 core/preflight/nightly_planner.py
ExecStart=%h/seevar/.venv/bin/python3 core/preflight/schedule_compiler.py
StandardOutput=append:%h/seevar/logs/planner.log
StandardError=append:%h/seevar/logs/planner.log

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ReadWritePaths=%h/seevar/data %h/seevar/logs
TimeoutStartSec=180s
TimeoutStopSec=10s

[Install]
WantedBy=default.target
